# installation-scripts/nxp-s32g/scripts/dreamos-setup.service
# This file is generated by k3s-master-prepare.sh
# It should be placed in /etc/systemd/system/ on the worker node.
[Unit]
Description=DreamOS Setup Service
After=network.target
Before=k3s.service
Wants=network.target

[Service]
Type=oneshot
ExecStart=/home/root/.dk/nxp-s32g/scripts/dreamos_setup.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
